# maj_db/crontab
0 1 * * * /usr/bin/python3 /async_scrap.py > /proc/1/fd/1 2>/proc/1/fd/2

